# jekyll_start.conf
# jekyll
start on startup

script
    cd /www/behaviac
    jekyll build --watch --config _config.yml
end script